# files/edgebox.list
system:system 755 /home/system/.edgebox_motd
system:system 757 /home/system/edgebox.sh
system:system 757 /home/system/tinc_boot
system:system 757 /home/system/edgeboxctl
root:root 757 /lib/systemd/system/edgeboxctl.service
